# build/mdta.conf
[program:mdta]
environment=DJANGO_SETTINGS_MODULE='mdta.settings.qaci01',PATH='/usr/local/virtualenvs/mdta/',HOME='/home/wicqacip'
directory=/opt/mdta
command = /usr/local/virtualenvs/mdta/bin/uwsgi --master --chdir /opt/mdta/ --module mdta.wsgi:application --socket /opt/mdta/run/mdta.sock --chmod-socket=666 --processes 3 --virtualenv /usr/local/virtualenvs/mdta/ --pythonpath /usr/local/virtualenvs/mdta/ --buffer-size 4096 --env DJANGO_SETTINGS_MODULE='mdta.settings.qaci01'

user = wicqacip
group = wicqacip
autostart=true
autorestart=true
stopsignal=QUIT
stderr_logfile = /var/log/mdta/mdta.log
stdout_logfile = /var/log/mdta/mdta.log
redirect_stderr = true

[program:mdta-celery]
environment=DJANGO_SETTINGS_MODULE='mdta.settings.qaci01',PATH='/usr/local/virtualenvs/mdta/'
directory=/opt/mdta/
command=/usr/local/virtualenvs/mdta/bin/celery worker -n qaci01.mdta  --app=mdta -Q mdta_queue --logfile=/var/log/mdta/mdta-celery.log -l warn -c 3
user=wicqacip
autostart=true
autorestart=true
redirect_stderr=True

